# wacom/.config/autostart/wacom-config.desktop
[Desktop Entry]
Type=Application
Terminal=false
Name=Wacom configs
Exec=/home/dzwiedz/scripts/wacom.sh
